# root/70-user.i3blocks.usb.rules
SUBSYSTEMS=="usb", RUN+="/usr/bin/pkill -RTMIN+3 i3blocks"
